# deploy/dedalo-ts.service
# Dédalo TS server — reference systemd unit (audit S2-38, DEC-17 item 3).
#
# Adjust the three ALL-CAPS placeholders (USER, /PATH/TO). The private config
# dir (../private relative to the repo) already holds .env; EnvironmentFile
# lets systemd overrides win over it (readEnv precedence: process env first).
#
#   sudo cp deploy/dedalo-ts.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now dedalo-ts
#
# Pair with dedalo-ts-watchdog.timer (health-checks /health over the unix
# socket and restarts on failure) — systemd's native WatchdogSec requires
# sd_notify, which Bun does not speak; the curl timer is the equivalent.

[Unit]
Description=Dedalo TS server (Bun)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=DEDALO_USER
WorkingDirectory=/PATH/TO/v7_ts/master_dedalo
# The runtime is PINNED (.bun-version / package.json engines.bun) — point
# ExecStart at the pinned binary, not a floating `bun` on PATH, so a stray
# `bun upgrade` cannot change the production runtime (audit S2-36).
ExecStart=/PATH/TO/.bun/bin/bun run src/server.ts
# Secrets/config: ../private/.env is read by the server itself; this line is
# for systemd-level OVERRIDES only (optional — remove the '-' file if unused).
EnvironmentFile=-/PATH/TO/v7_ts/private/.env

# Supervision (audit S2-38): any crash — OOM, native fault, unhandled error —
# restarts the server after 3 s. The server's own SIGTERM handler drains
# in-flight requests, closes the DB pool and unlinks the socket (S2-17).
Restart=always
RestartSec=3
TimeoutStopSec=30
KillSignal=SIGTERM

# Log capture (audit S2-37): stdout/stderr → journald. The structured access
# log (DEDALO_ACCESS_LOG=true) emits one JSON line per request here.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dedalo-ts

# Resource hygiene (tune to the host; the server is one shared process).
LimitNOFILE=65536
# MemoryMax=4G

[Install]
WantedBy=multi-user.target
