# deploy/dedalo-ts.service
# Dédalo TS server — reference systemd unit (audit S2-38, DEC-17 item 3).
#
# Adjust the ALL-CAPS/example placeholders below (DEDALO_USER, /opt/dedalo/…)
# to the real host layout. The canonical layout (see engineering/PRODUCTION.md)
# is:
#   /opt/dedalo/master_dedalo/   the repo (WorkingDirectory)
#   /opt/dedalo/private/         ../private (config, secrets, session store)
#   /opt/dedalo/.bun/bin/bun     the PINNED runtime binary
#   user/group: dedalo
#
#   sudo cp deploy/dedalo-ts.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now dedalo-ts
#
# Pair with dedalo-ts-watchdog.timer (health-checks /health over the unix
# socket and restarts on failure) — systemd's native WatchdogSec requires
# sd_notify, which Bun does not speak; the curl timer is the equivalent.

[Unit]
Description=Dedalo TS server (Bun)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=DEDALO_USER
WorkingDirectory=/opt/dedalo/master_dedalo
# The runtime is PINNED (.bun-version / package.json engines.bun) — point
# ExecStart at the pinned binary, not a floating `bun` on PATH, so a stray
# `bun upgrade` cannot change the production runtime (audit S2-36).
ExecStart=/opt/dedalo/.bun/bin/bun run src/server.ts

# Socket permissions (the 502-on-every-request trap). The server creates
# SERVER_UNIX_SOCKET under the process umask; the default 022 yields srwxr-xr-x, and
# CONNECTING to a unix socket requires WRITE permission — so a reverse proxy running as
# www-data/nginx gets EACCES and every request 502s, with nothing wrong in the app.
# UMask=0007 makes it srwxrwx---; add the web-server user to the `dedalo` group:
#     sudo usermod -aG dedalo www-data      # or: nginx
UMask=0007
# No EnvironmentFile here on purpose: the server reads ../private/.env
# ITSELF (src/config/env.ts readEnv/parseEnvFile), so pointing systemd's
# EnvironmentFile= at the same file is redundant — and actively harmful,
# because systemd's EnvironmentFile parser quotes/escapes differently from
# parseEnvFile and mangles raw-JSON values such as
# DEDALO_APPLICATION_LANGS={"lg-spa":"Castellano"} (the embedded unescaped
# quotes are not valid systemd env-file syntax). If a systemd-level override
# is genuinely needed, use `systemctl edit dedalo-ts` to add a drop-in
# Environment=/EnvironmentFile= instead of editing this unit.

# Supervision (audit S2-38): any crash — OOM, native fault, unhandled error —
# restarts the server after 3 s. The server's own SIGTERM handler drains
# in-flight requests, closes the DB pool and unlinks the socket (S2-17).
Restart=always
RestartSec=3
# Exit 75 (RESTART_EXIT_CODE, src/core/install/restart.ts) is a PLANNED restart:
# the installer persisted a new ../private/.env and the frozen-at-import config
# can only be reloaded by a fresh process. Restart=always already brings it back;
# this line just keeps the planned exit out of the failure counters (and out of
# the watchdog's alerting), instead of logging it as a crash.
SuccessExitStatus=75
TimeoutStopSec=30
KillSignal=SIGTERM

# Log capture (audit S2-37): stdout/stderr → journald. The structured access
# log (DEDALO_ACCESS_LOG=true) emits one JSON line per request here.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dedalo-ts

# Resource hygiene (tune to the host; the server is one shared process).
LimitNOFILE=65536
# MemoryMax=4G

[Install]
WantedBy=multi-user.target
